# etc/cron.d/StockMonitoring
*/1 *     * * *     root  cd /root/StockMonitoring/ && /home/python/bin/python script.py
